# contrib/docker/webapp/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=false

[program:apache2]
command=bash -c "service apache2 start"
