# etc/crontab
MAILTO=zack

# m h         dom mon dow user           command
1   0,6,12,18 *   *   *   debian-mirror  $DEBSRC/bin/update-debsources
