# contrib/docker/etc/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/bin/bash -c "service apache2 start"

#[programm:flask]
#command=/bin/bash -c "/opt/debsources/python/run_webapp.py --host 0.0.0.0"

[program:postgresql]
command=/bin/bash -c "service postgresql start"

[program:update-debsources]
command=/bin/bash -c "/opt/debsources/bin/update-debsources"
